# deployment/mcp.service
# =============================================================================
#   THE JESUS WEBSITE — SYSTEMD MCP SERVER SERVICE
#   File:    deployment/mcp.service
#   Version: 1.2.0
#   Purpose: Keeps the Read-Only FastMCP server running persistently
#            via stdio transport (communicates through stdin/stdout
#            rather than HTTP/port binding).
# =============================================================================

[Unit]
Description=The Jesus Website - MCP Read-Only Agent Interface (stdio transport)
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/thejesuswebsite
Environment="PATH=/var/www/thejesuswebsite/venv/bin"

# Launches mcp_server.py with default stdio transport — no port binding.
# External AI agents connect via stdin/stdout through the MCP protocol.
ExecStart=/var/www/thejesuswebsite/venv/bin/python mcp_server.py

# Standard auto-restarting pipeline
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
