# deployment/thejesuswebsite.service
# =============================================================================
#   THE JESUS WEBSITE — SYSTEMD MAIN SERVICE
#   File:    deployment/thejesuswebsite.service
#   Version: 1.1.0
#   Purpose: Keeps the FastAPI application running persistently via Unix socket.
# =============================================================================

[Unit]
Description=The Jesus Website - Main FastAPI Application
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/thejesuswebsite
Environment="PATH=/var/www/thejesuswebsite/venv/bin"

# Serves the full site via a Unix socket (proxied by nginx)
ExecStart=/var/www/thejesuswebsite/venv/bin/uvicorn serve_all:app --uds /var/www/thejesuswebsite/thejesuswebsite.sock --workers 2

Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
